# update-motd.d/80-ssl
#!/bin/sh
sudo certbot certificates | grep -B 1 'INVALID'
